# launchd/com.kezoo.pr-review-poller.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kezoo.pr-review-poller</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kezoo/.local/bin/pr-review-poller</string>
        <string>run</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kezoo/worktrees/.pr-review-poller-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kezoo/worktrees/.pr-review-poller-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/kezoo/.local/bin</string>
        <!-- Claude Code's HTTP-MCP OAuth lookup needs $USER to locate the
             cached token; without it plugin:slack:slack, github, datadog,
             notion all fall back to "Needs authentication" under launchd. -->
        <key>USER</key>
        <string>kezoo</string>
    </dict>
</dict>
</plist>
